# modules/rhizo_base/files/etc/sv/rapi/run
#!/bin/sh
export PYTHONPATH=$PYTHONPATH:/var/rhizomatica/rccn
#cd /var/rhizomatica/rccn ; python rapi.py debug
python /var/rhizomatica/rccn/rapi.py debug
