# modules/rhizo_base/files/etc/sv/freeswitch/run
#!/bin/sh

export PYTHONPATH=$PYTHONPATH:/var/rhizomatica/rccn
exec 2>&1
exec /usr/bin/freeswitch -nonat
