# modules/rhizo_base/files/etc/sv/kiwi/run
#!/bin/sh
exec /usr/bin/chpst -urhizomatica /usr/bin/nodejs /usr/share/kiwiirc/server/kiwi.js -f 
