# modules/rhizo_base/files/etc/sv/smpp/run
#!/bin/sh
export PYTHONPATH=$PYTHONPATH:/var/rhizomatica/rccn
exec python /var/rhizomatica/rccn/smpp.py
