# modules/rhizo_base/files/etc/sv/meas-json/run
#!/bin/bash
exec 2>&1

exec /usr/bin/stdbuf -oL /usr/bin/meas_json > /tmp/json_pipe
